# agent_server_no_llm_test/mcp_agent_isolated/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

# ── SAJHA MCP Server ──────────────────────────────────────────────────────────
# Starts first (priority 1). Agent server depends on SAJHA being up.
[program:sajha]
command=python run_sajha.py
directory=/app
environment=SAJHA_HOST="127.0.0.1",SAJHA_PORT="3002",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=1
stdout_logfile=/var/log/sajha.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/sajha_err.log
stderr_logfile_maxbytes=5MB

# ── Agent Server (FastAPI / uvicorn) ──────────────────────────────────────────
# Starts second (priority 2). nginx depends on this being up.
[program:agent]
command=uvicorn agent_server:app --host 127.0.0.1 --port 8000 --workers 1
directory=/app
environment=PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=2
stdout_logfile=/var/log/agent.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/agent_err.log
stderr_logfile_maxbytes=5MB

# ── nginx (reverse proxy + static file server) ───────────────────────────────
# Starts last (priority 3). Managed by the start.sh which preprocesses PORT.
[program:nginx]
command=nginx -g "daemon off;"
directory=/app
autostart=true
autorestart=true
startsecs=2
startretries=3
priority=3
stdout_logfile=/var/log/nginx_access.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/nginx_err.log
stderr_logfile_maxbytes=5MB
